# luci-qos-emong/files/etc/init.d/qos-emong
#!/bin/sh /etc/rc.common

START=99

start(){
/usr/bin/qos-emong start
}

stop(){
/usr/bin/qos-emong stop
}
restart(){
/usr/bin/qos-emong stop
/usr/bin/qos-emong start
}



